# scripts/cmpf_worker@.service
[Unit]
Description=CMPF worker instance %i
# 关联到 Worker 批量 Target
WantedBy=cmpf-workers.target

[Service]
Type=simple
ExecStart=/path/to/worker --worker-id=%i
Restart=on-failure
RestartSec=2
CPUQuota=20%
MemoryLimit=256M
User=cmpf
Group=cmpf
StandardOutput=journal
StandardError=journal